# .config/systemd/user/batsignal.service.d/options.conf
[Service]
ExecStart=
ExecStart=batsignal -w9 -c0 -d0 -f80 -NM "dunstify -u critical '%%s'"

